# sway/kanshi/kanshi.service
[Unit]
Description=kanshi output management

[Service]
Type=simple
ExecStart=kanshi

[Install]
WantedBy=graphical-session.target
