# sway/mako/mako.service
[Unit]
Description=mako notifications

[Service]
Type=simple
ExecStart=mako

[Install]
WantedBy=sway-session.target
