# sway/waybar/waybar.service
[Unit]
Description=Waybar

[Service]
Type=simple
ExecStart=waybar

[Install]
WantedBy=graphical-session.target
