# sway/xfce-polkit/xfce-polkit.service
[Unit]
Description=XFCE  Polkit

[Service]
Type=simple
ExecStart=/usr/libexec/xfce-polkit

[Install]
WantedBy=sway-session.target
